# ytarchive/rootfs/etc/services.d/ytarchive-run/run
#!/command/with-contenv bashio
# ==============================================================================
# Service to start YouTube livestream recording using ytarchive
# ==============================================================================

exec /usr/bin/ytarchive_start.sh