# ytarchive/rootfs/etc/services.d/ytarchive-stop/run
#!/command/with-contenv bashio
# ==============================================================================
# Service to stop YouTube livestream recording
# ==============================================================================

exec /usr/bin/ytarchive_stop.sh

